# netmaster.service.example
[Unit]
Description=netmaster CEC TV control
After=network.target

[Service]
ExecStart=/usr/bin/python3 /home/YOUR_USER/netmaster/server.py
WorkingDirectory=/home/YOUR_USER/netmaster
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
