# tvmaster.service.example
[Unit]
Description=tvmaster CEC TV control
After=network.target

[Service]
ExecStart=/home/YOUR_USER/tvmaster/venv/bin/python /home/YOUR_USER/tvmaster/server.py
WorkingDirectory=/home/YOUR_USER/tvmaster
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
